# sandbox/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; Path of UNIX socket file
chmod=0770                  ; Socket file permissions
chown=ubuntu:ubuntu         ; Socket file owner (matches USER ubuntu)
username=%(ENV_SUPERVISOR_RPC_USERNAME)s ; RPC username
password=%(ENV_SUPERVISOR_RPC_PASSWORD)s ; RPC password

[supervisord]
logfile=/dev/null           ; Suppress logfile (nodaemon=true already writes to stdout)
logfile_maxbytes=0          ; No limit on log size
loglevel=info               ; Log level
pidfile=/tmp/supervisord.pid ; Path of pidfile
nodaemon=true                ; Need to run in foreground in Docker
minfds=1024                  ; Minimum value of open file descriptors
minprocs=200                 ; Minimum number of processes
autoshutdown=true            ; Auto exit after all services stop

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; Socket connection to supervisord
username=%(ENV_SUPERVISOR_RPC_USERNAME)s ; RPC username
password=%(ENV_SUPERVISOR_RPC_PASSWORD)s ; RPC password

; Runtime dirs for ubuntu user (runs as ubuntu, no root needed)
; /run/user/1000 provided by tmpfs mount in compose; /tmp and .cache writable
[program:runtime_init]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu/dconf /home/ubuntu/.config /tmp/.X11-unix 2>/dev/null; chmod 1777 /tmp/.X11-unix 2>/dev/null; exit 0"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=1

; GitHub CLI authentication (one-shot, runs after runtime_init)
[program:gh_auth_setup]
command=/bin/bash -c "if [ -n \"%(ENV_GH_TOKEN)s\" ] && [ \"%(ENV_GH_TOKEN)s\" != \"\" ]; then echo '%(ENV_GH_TOKEN)s' | gh auth login --with-token 2>/dev/null && gh auth setup-git 2>/dev/null && echo 'gh auth: configured'; else echo 'gh auth: skipped (no GH_TOKEN)'; fi; exit 0"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,GH_TOKEN=%(ENV_GH_TOKEN)s
user=ubuntu
priority=3

; Environment context generator - runs once at startup
[program:context_generator]
command=python3 /app/scripts/generate_sandbox_context.py
directory=/app
user=ubuntu
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,SANDBOX_CONTEXT_JSON=/app/sandbox_context.json,SANDBOX_CONTEXT_MD=/app/sandbox_context.md
priority=5

; Local session D-Bus to avoid repeated Chrome bus connection errors in minimal containers
[program:dbus_session]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu; exec dbus-daemon --session --address=unix:path=/tmp/runtime-ubuntu/dbus-session --nofork --nopidfile --nosyslog"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu
user=ubuntu
priority=8

; ============================================================================
; Display Stack (Xvfb + Window Manager + VNC)
; ============================================================================

; Virtual framebuffer — provides X11 display :99 for headed Chrome
[program:xvfb]
command=Xvfb :99 -screen 0 1280x1024x24 -ac -nolisten tcp -dpi 96 +extension GLX +extension RANDR +render -noreset
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
user=ubuntu
priority=10

; Openbox window manager — minimal WM for Chrome window management
[program:openbox]
command=/bin/bash -c "while ! xdpyinfo -display :99 >/dev/null 2>&1; do sleep 0.2; done; exec openbox --sm-disable --config-file /app/scripts/openbox-rc.xml"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,DISPLAY=:99
user=ubuntu
priority=15

; x11vnc — VNC server attached to Xvfb display
[program:x11vnc]
command=/bin/bash -c "while ! xdpyinfo -display :99 >/dev/null 2>&1; do sleep 0.2; done; exec x11vnc -display :99 -forever -nopw -shared -rfbport 5900 -xkb -noxrecord -noxfixes -noxdamage -no6"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,DISPLAY=:99
user=ubuntu
priority=20

; Websockify — bridges VNC (TCP 5900) to WebSocket (TCP 6080) for noVNC
[program:websockify]
command=/bin/bash -c "exec websockify --heartbeat=30 0.0.0.0:6080 localhost:5900 2>&1 | grep -v 'no.*numpy'"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
user=ubuntu
priority=25

; ============================================================================
; Chrome (headed on Xvfb display :99)
; ============================================================================

; Chrome in headed mode on virtual display — CDP screencast and VNC both work
; Uses wrapper script for clean stderr filtering (no orphan processes)
[program:chrome_cdp_only]
command=/bin/bash /app/scripts/run_chrome.sh
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,DISPLAY=:99,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu,BROWSER_PATH=%(ENV_BROWSER_PATH)s
user=ubuntu
priority=50
startretries=5
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; socat port forwarding configuration
; Always starts (CDP port forwarding works in both modes)
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=30
startsecs=2

; ============================================================================
; Core Services - Always run in both modes
; ============================================================================

; FastAPI application configuration
[program:app]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=50

; Sandbox framework API (agent DB + workflow coordination)
[program:framework]
command=python3 -m uvicorn app.framework.main:app --host 0.0.0.0 --port 8082 %(ENV_FRAMEWORK_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=55

; ============================================================================
; code-server — VS Code web IDE (optional, enabled by ENABLE_CODE_SERVER env var)
; ============================================================================
[program:code_server]
command=/bin/bash -c "if [ '%(ENV_ENABLE_CODE_SERVER)s' = '1' ] || [ '%(ENV_ENABLE_CODE_SERVER)s' = 'true' ]; then exec code-server --bind-addr 0.0.0.0:%(ENV_CODE_SERVER_PORT)s --auth password --disable-telemetry --disable-update-check /workspace; else echo 'code-server disabled'; exit 0; fi"
autostart=%(ENV_ENABLE_CODE_SERVER)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/ubuntu",PASSWORD="%(ENV_CODE_SERVER_PASSWORD)s"
user=ubuntu
priority=60

; ============================================================================
; Group configuration
; CDP-only service group
[group:services]
programs=runtime_init,gh_auth_setup,context_generator,dbus_session,xvfb,openbox,x11vnc,websockify,chrome_cdp_only,socat,app,framework,code_server
