# sandbox/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; Path of UNIX socket file
chmod=0770                  ; Socket file permissions
chown=ubuntu:ubuntu         ; Socket file owner (matches USER ubuntu)

[supervisord]
logfile=/dev/null           ; Suppress logfile (nodaemon=true already writes to stdout)
logfile_maxbytes=0          ; No limit on log size
loglevel=info               ; Log level
pidfile=/tmp/supervisord.pid ; Path of pidfile
nodaemon=true                ; Need to run in foreground in Docker
minfds=1024                  ; Minimum value of open file descriptors
minprocs=200                 ; Minimum number of processes
autoshutdown=true            ; Auto exit after all services stop

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; Socket connection to supervisord

; Runtime dirs for ubuntu user (runs as ubuntu, no root needed)
; /run/user/1000 provided by tmpfs mount in compose; /tmp and .cache writable
[program:runtime_init]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu/dconf /home/ubuntu/.cache/openbox /home/ubuntu/.config 2>/dev/null; exit 0"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=1

; Environment context generator - runs once at startup
[program:context_generator]
command=python3 /app/scripts/generate_sandbox_context.py
directory=/app
user=ubuntu
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,SANDBOX_CONTEXT_JSON=/app/sandbox_context.json,SANDBOX_CONTEXT_MD=/app/sandbox_context.md
priority=5

; Local session D-Bus to avoid repeated Chrome bus connection errors in minimal containers
[program:dbus_session]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu; exec dbus-daemon --session --address=unix:path=/tmp/runtime-ubuntu/dbus-session --nofork --nopidfile --nosyslog"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu
user=ubuntu
priority=8

; ============================================================================
; CDP-Only Mode: Chrome runs in headless mode, no X11/VNC stack
; Dual Mode: Chrome runs on Xvfb display with VNC access
; Controlled by SANDBOX_STREAMING_MODE environment variable
; ============================================================================

; Chrome in CDP-only mode (headless, no DISPLAY)
; Only starts when SANDBOX_STREAMING_MODE=cdp_only
[program:chrome_cdp_only]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" != \"cdp_only\" ]; then echo 'Skipping chrome_cdp_only (mode: ${SANDBOX_STREAMING_MODE:-dual})'; exit 0; fi; while [ ! -S /tmp/runtime-ubuntu/dbus-session ]; do echo 'Waiting for dbus session...'; sleep 0.2; done; mkdir -p /tmp/runtime-ubuntu /tmp/runtime-ubuntu/dconf; export XDG_RUNTIME_DIR=/tmp/runtime-ubuntu; export DBUS_SESSION_BUS_ADDRESS=unix:path=/tmp/runtime-ubuntu/dbus-session; export DBUS_SYSTEM_BUS_ADDRESS=unix:path=/tmp/runtime-ubuntu/dbus-session; export GSETTINGS_BACKEND=memory; sleep 1; exec \"${BROWSER_PATH:-/opt/chrome-for-testing/chrome}\" --headless=new --window-size=1280,1024 --disable-dev-shm-usage --disable-gpu --disable-gpu-compositing --disable-gpu-sandbox --disable-vulkan --use-angle=swiftshader --use-gl=swiftshader --enable-unsafe-swiftshader --disable-accelerated-2d-canvas --disable-accelerated-jpeg-decoding --disable-accelerated-mjpeg-decode --disable-accelerated-video-decode --no-first-run --no-default-browser-check --test-type --noerrdialogs --disable-session-crashed-bubble --hide-crash-restore-bubble --disable-features=WelcomeExperience,SigninPromo,TranslateUI,BlinkGenPropertyTrees,AudioServiceOutOfProcess,InfiniteSessionRestore,VizDisplayCompositor,GCMChannelStatusRequest,MediaRouter,DialMediaRouteProvider,PushMessaging,OptimizationHints,AutofillServerCommunication --enable-features=NetworkService,NetworkServiceInProcess --disable-infobars --disable-notifications --disable-popup-blocking --disable-prompt-on-repost --disable-extensions --disable-component-extensions-with-background-pages --disable-component-update --disable-component-cloud-policy --disable-background-networking --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-breakpad --disable-client-side-phishing-detection --disable-default-apps --disable-domain-reliability --disable-hang-monitor --disable-ipc-flooding-protection --disable-renderer-backgrounding --disable-sync --disable-translate --metrics-recording-only --mute-audio --no-pings --password-store=basic --use-mock-keychain --force-color-profile=srgb --force-device-scale-factor=1 --disable-smooth-scrolling --disable-threaded-scrolling --disable-blink-features=AutomationControlled --user-agent='Mozilla/5.0 (X11; Linux x86_64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/128.0.6613.137 Safari/537.36' --lang=en-US --remote-debugging-address=0.0.0.0 --remote-debugging-port=8222 --disable-gcm-client %(ENV_CHROME_ARGS)s 2> >(python3 /app/scripts/chrome_stderr_filter.py >&2)"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu,BROWSER_PATH=%(ENV_BROWSER_PATH)s,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=50
startretries=5
startsecs=5
stopwaitsecs=10

; Chrome in dual mode (X11 display with VNC)
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:chrome_dual]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping chrome_dual (mode: cdp_only)'; exit 0; fi; while ! xdpyinfo -display :1 >/dev/null 2>&1; do echo 'Waiting for display :1...'; sleep 1; done; while ! pgrep -x openbox >/dev/null 2>&1; do echo 'Waiting for openbox...'; sleep 0.5; done; while [ ! -S /tmp/runtime-ubuntu/dbus-session ]; do echo 'Waiting for dbus session...'; sleep 0.2; done; mkdir -p /tmp/runtime-ubuntu /tmp/runtime-ubuntu/dconf; export XDG_RUNTIME_DIR=/tmp/runtime-ubuntu; export DBUS_SESSION_BUS_ADDRESS=unix:path=/tmp/runtime-ubuntu/dbus-session; export DBUS_SYSTEM_BUS_ADDRESS=unix:path=/tmp/runtime-ubuntu/dbus-session; export GSETTINGS_BACKEND=memory; sleep 1; exec \"${BROWSER_PATH:-/opt/chrome-for-testing/chrome}\" --display=:1 --window-position=0,0 --window-size=1280,1024 --disable-dev-shm-usage --disable-gpu --disable-gpu-compositing --disable-gpu-sandbox --disable-vulkan --use-angle=swiftshader --use-gl=swiftshader --enable-unsafe-swiftshader --disable-accelerated-2d-canvas --disable-accelerated-jpeg-decoding --disable-accelerated-mjpeg-decode --disable-accelerated-video-decode --no-first-run --no-default-browser-check --no-xshm --test-type --noerrdialogs --disable-session-crashed-bubble --hide-crash-restore-bubble --disable-features=WelcomeExperience,SigninPromo,TranslateUI,BlinkGenPropertyTrees,AudioServiceOutOfProcess,InfiniteSessionRestore,VizDisplayCompositor,GCMChannelStatusRequest,MediaRouter,DialMediaRouteProvider,PushMessaging,OptimizationHints,AutofillServerCommunication --enable-features=NetworkService,NetworkServiceInProcess --disable-infobars --disable-notifications --disable-popup-blocking --disable-prompt-on-repost --disable-extensions --disable-component-extensions-with-background-pages --disable-component-update --disable-component-cloud-policy --disable-background-networking --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-breakpad --disable-client-side-phishing-detection --disable-default-apps --disable-domain-reliability --disable-hang-monitor --disable-ipc-flooding-protection --disable-renderer-backgrounding --disable-sync --disable-translate --metrics-recording-only --mute-audio --no-pings --password-store=basic --use-mock-keychain --force-color-profile=srgb --force-device-scale-factor=1 --disable-smooth-scrolling --disable-threaded-scrolling --disable-blink-features=AutomationControlled --user-agent='Mozilla/5.0 (X11; Linux x86_64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/128.0.6613.137 Safari/537.36' --lang=en-US --remote-debugging-address=0.0.0.0 --remote-debugging-port=8222 --disable-gcm-client %(ENV_CHROME_ARGS)s 2> >(python3 /app/scripts/chrome_stderr_filter.py >&2)"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu,BROWSER_PATH=%(ENV_BROWSER_PATH)s,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=50
startretries=5
startsecs=5
stopwaitsecs=10

; socat port forwarding configuration
; Always starts (CDP port forwarding works in both modes)
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=30
startsecs=2

; ============================================================================
; X11/VNC Stack - Only runs in dual mode
; ============================================================================

; Xvfb - Virtual display server
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:xvfb]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping Xvfb (mode: cdp_only)'; exit 0; fi; exec /usr/bin/Xvfb :1 -screen 0 1280x1024x24 -ac +extension GLX +extension RANDR +render -noreset -nolisten tcp -dpi 96 -shmem"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,USER=ubuntu,XDG_RUNTIME_DIR=/run/user/1000,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=10
startsecs=3
stopwaitsecs=5

; Openbox - Minimal window manager
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:openbox]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping Openbox (mode: cdp_only)'; exit 0; fi; while ! xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.5; done; exec openbox"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=12
startsecs=2
stopwaitsecs=5

; xrandr setup - CRITICAL for fixing VNC offset issue
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:xrandr_setup]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping xrandr (mode: cdp_only)'; exit 0; fi; while ! xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.5; done; sleep 1; xrandr --display :1 --fb 1280x1024; echo 'xrandr framebuffer set to 1280x1024'"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=13

; x11vnc - VNC server
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:x11vnc]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping x11vnc (mode: cdp_only)'; exit 0; fi; exec /usr/bin/x11vnc -display :1 -forever -shared -nopw -rfbport 5900 -no6 -xkb -defer 10 -wait 1 -threads -allinput -cursor arrow -sb 0 -noxrecord -noxdamage -noxinerama -clip 1280x1024+0+0"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,USER=ubuntu,XDG_RUNTIME_DIR=/run/user/1000,SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=15
startsecs=2
stopwaitsecs=5

; Websockify service - WebSocket to VNC proxy
; Only starts when SANDBOX_STREAMING_MODE=dual (default)
[program:websockify]
command=/bin/bash -c "if [ \"${SANDBOX_STREAMING_MODE:-dual}\" = \"cdp_only\" ]; then echo 'Skipping websockify (mode: cdp_only)'; exit 0; fi; exec python3 -m websockify 0.0.0.0:5901 localhost:5900"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SANDBOX_STREAMING_MODE=%(ENV_SANDBOX_STREAMING_MODE)s
user=ubuntu
priority=20
startsecs=2
stopwaitsecs=5

; ============================================================================
; Core Services - Always run in both modes
; ============================================================================

; FastAPI application configuration
[program:app]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=50

; Sandbox framework API (agent DB + workflow coordination)
[program:framework]
command=python3 -m uvicorn app.framework.main:app --host 0.0.0.0 --port 8082 %(ENV_FRAMEWORK_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=55

; Group configuration
; Note: chrome_cdp_only and chrome_dual are mutually exclusive based on env var
[group:services]
programs=runtime_init,context_generator,dbus_session,xvfb,openbox,xrandr_setup,x11vnc,websockify,chrome_cdp_only,chrome_dual,socat,app,framework
