# sandbox/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; Path of UNIX socket file, supervisor uses it to listen
chmod=0770                  ; Socket file permissions
chown=ubuntu:ubuntu         ; Socket file owner and group

[supervisord]
logfile=/dev/null           ; Suppress logfile (nodaemon=true already writes to stdout)
logfile_maxbytes=0          ; No limit on log size
loglevel=info               ; Log level
pidfile=/tmp/supervisord.pid ; Path of pidfile
nodaemon=true                ; Need to run in foreground in Docker
minfds=1024                  ; Minimum value of open file descriptors
minprocs=200                 ; Minimum number of processes
autoshutdown=true            ; Auto exit after all services stop

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; Socket connection to supervisord

; Fix tmpfs permissions - runs once at startup before everything else
[program:fix_permissions]
command=/bin/bash /app/fix-permissions.sh
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
priority=1

; Environment context generator - runs once at startup
[program:context_generator]
command=python3 /app/scripts/generate_sandbox_context.py
directory=/app
user=ubuntu
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,SANDBOX_CONTEXT_JSON=/app/sandbox_context.json,SANDBOX_CONTEXT_MD=/app/sandbox_context.md
priority=5

; Note: Xvfb provides the virtual display, x11vnc provides VNC access
; Chrome for Testing 128.0.6613.137 - waits for display :1 and openbox before starting
; Using explicit window position/size without --start-maximized to avoid conflicts
[program:chrome]
command=/bin/bash -c "while ! xdpyinfo -display :1 >/dev/null 2>&1; do echo 'Waiting for display :1...'; sleep 1; done; while ! pgrep -x openbox >/dev/null 2>&1; do echo 'Waiting for openbox...'; sleep 0.5; done; mkdir -p /tmp/runtime-ubuntu /tmp/runtime-ubuntu/dconf; export XDG_RUNTIME_DIR=/tmp/runtime-ubuntu; sleep 1; exec ${CHROME_BIN:-/opt/chrome-for-testing/chrome} --display=:1 --window-position=0,0 --window-size=1280,1024 --disable-dev-shm-usage --disable-gpu --disable-gpu-compositing --disable-gpu-sandbox --disable-vulkan --use-angle=swiftshader --use-gl=swiftshader --enable-unsafe-swiftshader --disable-accelerated-2d-canvas --disable-accelerated-jpeg-decoding --disable-accelerated-mjpeg-decode --disable-accelerated-video-decode --no-first-run --no-default-browser-check --no-xshm --test-type --noerrdialogs --disable-session-crashed-bubble --hide-crash-restore-bubble --disable-features=WelcomeExperience,SigninPromo,TranslateUI,BlinkGenPropertyTrees,AudioServiceOutOfProcess,InfiniteSessionRestore,VizDisplayCompositor,GCMChannelStatusRequest --enable-features=NetworkService,NetworkServiceInProcess --disable-infobars --disable-notifications --disable-popup-blocking --disable-prompt-on-repost --disable-extensions --disable-component-extensions-with-background-pages --disable-background-networking --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-breakpad --disable-client-side-phishing-detection --disable-default-apps --disable-domain-reliability --disable-hang-monitor --disable-ipc-flooding-protection --disable-renderer-backgrounding --disable-sync --disable-translate --metrics-recording-only --mute-audio --no-pings --password-store=basic --use-mock-keychain --force-color-profile=srgb --force-device-scale-factor=1 --disable-smooth-scrolling --disable-threaded-scrolling --disable-blink-features=AutomationControlled --user-agent='Mozilla/5.0 (X11; Linux x86_64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/128.0.6613.137 Safari/537.36' --lang=en-US --remote-debugging-address=0.0.0.0 --remote-debugging-port=8222 --dbus-stub %(ENV_CHROME_ARGS)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu,CHROME_BIN=/opt/chrome-for-testing/chrome
user=ubuntu
priority=50
startretries=5
startsecs=5
stopwaitsecs=10

; socat port forwarding configuration
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=30
startsecs=2

; Xvfb - Virtual display server
; Provides display :1 for Chrome and other X11 applications
; PERFORMANCE OPTIMIZATIONS:
; -nolisten tcp: Security - no TCP listening
; +extension RANDR: Enable resize/rotate extension
; -dpi 96: Consistent DPI for rendering
; -shmem: Use shared memory for better performance
[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x1024x24 -ac +extension GLX +extension RANDR +render -noreset -nolisten tcp -dpi 96 -shmem
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,USER=ubuntu,XDG_RUNTIME_DIR=/run/user/1000
user=ubuntu
priority=10
startsecs=3
stopwaitsecs=5

; Openbox - Minimal window manager for proper window management
; Required for Chrome to maximize properly on Xvfb
[program:openbox]
command=/bin/bash -c "while ! xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.5; done; exec openbox"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu
user=ubuntu
priority=12
startsecs=2
stopwaitsecs=5

; xrandr setup - CRITICAL for fixing VNC offset issue
; Forces framebuffer to correct size after Xvfb starts
; This prevents the geometry mismatch between Xvfb and x11vnc
[program:xrandr_setup]
command=/bin/bash -c "while ! xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.5; done; sleep 1; xrandr --display :1 --fb 1280x1024; echo 'xrandr framebuffer set to 1280x1024'"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1,HOME=/home/ubuntu
user=ubuntu
priority=13

; x11vnc - VNC server
; Provides VNC access to display :1 on port 5900
; PERFORMANCE OPTIMIZATIONS:
; -defer 10: 10ms update batching (reduced from 15ms for lower latency)
; -wait 1: Minimal initial wait time
; -threads: Multi-threaded encoding for better CPU utilization
; -allinput: Better keyboard/mouse input handling
; -cursor arrow: Lightweight cursor rendering
; -sb 0: No screen blank timeout
; -noxrecord: Disable XRECORD (reduces overhead)
; -clip: Clip to screen size to prevent any overflow
; NOTE: -ncache removed as it can cause display offset/artifacts
[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw -rfbport 5900 -xkb -defer 10 -wait 1 -threads -allinput -cursor arrow -sb 0 -noxrecord -noxdamage -noxinerama -clip 1280x1024+0+0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,USER=ubuntu,XDG_RUNTIME_DIR=/run/user/1000
user=ubuntu
priority=15
startsecs=2
stopwaitsecs=5

; Websockify service - WebSocket to VNC proxy
; Provides WebSocket interface on port 5901 for noVNC clients
; Proxies to native VNC on port 5900
[program:websockify]
command=python3 -m websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=20
startsecs=2
stopwaitsecs=5

; FastAPI application configuration
[program:app]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=50

; Sandbox framework API (agent DB + workflow coordination)
[program:framework]
command=uvicorn app.framework.main:app --host 0.0.0.0 --port 8082 %(ENV_FRAMEWORK_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=55

; Group configuration, can start or stop multiple programs at once
[group:services]
programs=fix_permissions,context_generator,xvfb,openbox,xrandr_setup,x11vnc,websockify,chrome,socat,app,framework
