# sandbox/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; Path of UNIX socket file
chmod=0770                  ; Socket file permissions
chown=ubuntu:ubuntu         ; Socket file owner (matches USER ubuntu)

[supervisord]
logfile=/dev/null           ; Suppress logfile (nodaemon=true already writes to stdout)
logfile_maxbytes=0          ; No limit on log size
loglevel=info               ; Log level
pidfile=/tmp/supervisord.pid ; Path of pidfile
nodaemon=true                ; Need to run in foreground in Docker
minfds=1024                  ; Minimum value of open file descriptors
minprocs=200                 ; Minimum number of processes
autoshutdown=true            ; Auto exit after all services stop

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; Socket connection to supervisord

; Runtime dirs for ubuntu user (runs as ubuntu, no root needed)
; /run/user/1000 provided by tmpfs mount in compose; /tmp and .cache writable
[program:runtime_init]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu/dconf /home/ubuntu/.config 2>/dev/null; exit 0"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=1

; Environment context generator - runs once at startup
[program:context_generator]
command=python3 /app/scripts/generate_sandbox_context.py
directory=/app
user=ubuntu
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,SANDBOX_CONTEXT_JSON=/app/sandbox_context.json,SANDBOX_CONTEXT_MD=/app/sandbox_context.md
priority=5

; Local session D-Bus to avoid repeated Chrome bus connection errors in minimal containers
[program:dbus_session]
command=/bin/bash -c "mkdir -p /tmp/runtime-ubuntu; exec dbus-daemon --session --address=unix:path=/tmp/runtime-ubuntu/dbus-session --nofork --nopidfile --nosyslog"
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu
user=ubuntu
priority=8

; ============================================================================
; Chrome (CDP-only, headless)
; ============================================================================

; Chrome in CDP-only mode (headless, no DISPLAY)
; Uses wrapper script for clean stderr filtering (no orphan processes)
[program:chrome_cdp_only]
command=/bin/bash /app/scripts/run_chrome.sh
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,XDG_RUNTIME_DIR=/tmp/runtime-ubuntu,BROWSER_PATH=%(ENV_BROWSER_PATH)s
user=ubuntu
priority=50
startretries=5
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; socat port forwarding configuration
; Always starts (CDP port forwarding works in both modes)
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=ubuntu
priority=30
startsecs=2

; ============================================================================
; Core Services - Always run in both modes
; ============================================================================

; FastAPI application configuration
[program:app]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=50

; Sandbox framework API (agent DB + workflow coordination)
[program:framework]
command=python3 -m uvicorn app.framework.main:app --host 0.0.0.0 --port 8082 %(ENV_FRAMEWORK_UVI_ARGS)s
directory=/app
user=ubuntu
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=55

; ============================================================================
; Optional Addon Services (SANDBOX_ADDONS_ENABLED=1)
; ============================================================================

; code-server 4.104.3 — VS Code in browser on port 8443 (addons only)
; Exits cleanly when addons are not enabled; supervisord does not restart it.
[program:code_server]
command=/bin/bash -c 'if [ "${SANDBOX_ADDONS_ENABLED:-0}" = "1" ] && command -v code-server >/dev/null 2>&1; then exec code-server --bind-addr 0.0.0.0:8443 --auth none /workspace; fi'
autostart=true
autorestart=false
exitcodes=0
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu,PASSWORD=""
user=ubuntu
priority=60

; Group configuration
; CDP-only service group
[group:services]
programs=runtime_init,context_generator,dbus_session,chrome_cdp_only,socat,app,framework,code_server
